# 02-module/22-systemd/etc/systemd/system/custom.target
[Unit]
Description=Custom Target
Requires=multi-user.target
Wants=varsize.service sshd.service

[Install]
WantedBy=multi-user.target
