# 02-module/22-systemd/etc/systemd/system/my.service
[Unit]
Description=My Super service

[Service]
ExecStart=/root/my.sh
Restart=always

[Install]
WantedBy=multi-user.target
