# 02-module/21-cron/etc/cron.d/uptime
0 9 * * 1-5 root uptime >> /tmp/uptime.log
